# chef-repo/site-cookbooks/twonkymedia/files/default/twonkymedia.service
[Unit]
Description=twonkymedia

[Service]
ExecStart=/usr/local/twonky/twonkymedia.sh -D $OPTIONS

[Install]
WantedBy=multi-user.target